# deploy/ec2/synapse-frontend.service
[Unit]
Description=Synapse Next.js frontend
After=network.target synapse-api.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Synapse/frontend
Environment=NODE_ENV=production
EnvironmentFile=-/home/ubuntu/Synapse/frontend/.env.production
ExecStart=/home/ubuntu/Synapse/frontend/node_modules/.bin/next start -H 127.0.0.1 -p 3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
